# resources/hindsight-rebuild.service.template
[Unit]
Description=Hindsight FAISS Index Rebuild
[Service]
ExecStart=EXECSTART_PLACEHOLDER