# resources/hindsight-rebuild.service.template
[Unit]
Description=Hindsight - Rebuild FAISS and Recoll Index

[Service]
Type=oneshot
ExecStart=%%VENV_PATH%%/bin/python %%APP_PATH%%/rebuild_index.py
